# roles/ko/nest/templates/nest.service
[Unit]
Description=Run nest provision

[Service]
ExecStart={{ nest_autorun_cmd }}
Type=oneshot

[Install]
WantedBy=default.target
