# roles/ko/nest/templates/systemd.service.j2
[Unit]
Description=Run nest provision

[Service]
WorkingDirectory={{ nest_path }}/
ExecStart=ansible-playbook -i inventory/nest site.yml
Type=oneshot

[Install]
WantedBy=default.target
